# roles/local-content-proxy/templates/content-proxy.service.j2
[Unit]
Description=Rustici Content Proxy
Requires=network-online.target
After=network-online.target

[Service]
EnvironmentFile={{ content_proxy_installation_dir }}/environment-variables
ExecStart=/usr/bin/java -jar "{{ content_proxy_installation_dir }}/content-proxy-{{ content_proxy_release }}.jar" -Dvertx.logger-delegate-factory-class-name=io.vertx.core.logging.SLF4JLogDelegateFactory

User=tomcat
Group=tomcat

[Install]
WantedBy=multi-user.target
