# roles/content-controller/templates/cchealthcheck.service.j2
[Unit]
Description=Content Controller Healthcheck Daemon
After=network.target

[Service]
ExecStart=/usr/local/bin/cc-healthcheck.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
